# aprigatto.service
[Unit]
Description=Apri Gatto
After=multi-user.target

[Install]
WantedBy=multi-user.target

[Service]
Type=simple
ExecStart=/home/robot/main